# ansible/roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/node_exporter \
    --web.listen-address="127.0.0.1:{{ node_exporter_port }}" \
    --collector.systemd \
    --collector.processes
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
